# templates/default/procfile.init.erb
#!/bin/bash
### BEGIN INIT INFO
# Provides:          selenium
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start selenium at boot time
# Description:       Enable service provided by selenium.
### END INIT INFO
 
  <% node[@name.to_sym].each do |k, v| %>
  export <%= k.to_s.upcase %>='<%= v.to_s %>'
  <% end %>
 
 
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
 
DISPLAY_ID=":99"
RUN_AS=root
 
JAVA_BIN=/usr/bin/java
 
SELENIUM_DIR=<%=node['selenium']['server']['installpath']=%>
SELENIUM_PID_FILE="$SELENIUM_DIR/selenium.pid"
SELENIUM_JAR_FILE="$SELENIUM_DIR/selenium-server-standalone.jar"
SELENIUM_LOG_FILE=<%=node['selenium']['server']['logpath']%>
SELENIUM_DAEMON_OPTS=" -client -jar $SELENIUM_JAR_FILE -log $SELENIUM_LOG_FILE"
 
export DISPLAY="$DISPLAY_ID.0"
 
set -e
 
. /lib/lsb/init-functions
 
case "$1" in
    start)
 
        if status_of_proc -p $SELENIUM_PID_FILE "$SELENIUM_JAR_FILE" $SELENIUM_JAR_FILE > /dev/null; then
            log_progress_msg "Service already running"
        else
            log_daemon_msg "Starting Selenium server"
            log_progress_msg "selenium"
            start-stop-daemon -c $RUN_AS --start --quiet --background --pidfile $SELENIUM_PID_FILE --make-pidfile --exec $JAVA_BIN -- $SELENIUM_DAEMON_OPTS
        fi
        ;;
 
    stop)
        if status_of_proc -p $SELENIUM_PID_FILE "$SELENIUM_JAR_FILE" $SELENIUM_JAR_FILE > /dev/null; then
            log_daemon_msg "Stopping Selenium server"
            log_progress_msg "selenium"
            start-stop-daemon --stop --pidfile $SELENIUM_PID_FILE 
        else
            log_progress_msg "Service not running"
        fi
        ;;
 
    restart|force-reload)
        $0 stop
        sleep 1
        $0 start
	;;
 
	status)
        status_of_proc -p $SELENIUM_PID_FILE "$SELENIUM_JAR_FILE" $SELENIUM_JAR_FILE && exit 0 || exit $?
    ;;
 
    *)
      	N=/etc/init.d/selenium
        echo "Usage: $N {start|stop|restart|force-reload|status}" >&2
        exit 1
        ;;
esac